# templates/default/autossh.service.erb
# Generated by Chef! Do not edit, your changes will be overriden!
# Cookbook Name:: autossh

[Unit]
Description=AutoSSH Service
After=network-online.target
Wants=network-online.target

[Service]
User=<%= @luser %>
Environment="AUTOSSH_GATETIME=0"
ExecStart=<%= @bin %> <%= @flags %> <%= @lport %>:<%= @bind %>:<%= @rport %> <%= @ruser %>@<%= @ssh_host %> -p <%= @ssh_port %>

[Install]
WantedBy=network-online.target
